# vcan.service
[Unit]
Description=vcan.service

[Service]
ExecStart=/root/vcan.sh

[Install]
WantedBy=multi-user.target
